# launchd/com.metacogna.claude.quarterly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.metacogna.claude.quarterly-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nullzero/Library/Application Support/reflex/bun/bin/bun</string>
        <string>run</string>
        <string>review:product:quarterly</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nullzero/Metacogna/claude-obsidian-todoist-linear</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/nullzero</string>
    </dict>

    <!-- Run on Jan 1, Apr 1, Jul 1, Oct 1 at 10:00 AM -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Month</key>
            <integer>1</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>4</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>7</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>10</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-quarterly-review.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-quarterly-review-error.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>

    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
