# launchd/com.metacogna.claude.daily-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.metacogna.claude.daily-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nullzero/Library/Application Support/reflex/bun/bin/bun</string>
        <string>run</string>
        <string>review:daily</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nullzero/Metacogna/claude-obsidian-todoist-linear</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/nullzero</string>
    </dict>

    <!-- Run Monday through Friday at 9:00 AM -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-daily-review.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-daily-review-error.log</string>

    <!-- Keep job loaded even if it exits -->
    <key>KeepAlive</key>
    <false/>

    <!-- Don't run immediately on load -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Allow job to run even when user is not logged in -->
    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
